# startup/daphne.conf
[program:daphne]
command=/var/www/gtracker/startup/daphne_start.sh
user=gtracker
stdout_logfile=/var/www/gtracker/logs/daphne.log
redirect_stderr=true
autostart=true
autorestart=true